# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:mongodb]
command=/usr/bin/mongod

[program:action-jackson]
command=/usr/bin/node /action-jackson-master/index.js



